# infra/atlas/ansible/roles/atlas_homeserver/templates/taskchampion-backup.service.j2
[Unit]
Description=Backup TaskChampion sync data to truenas

[Service]
Type=oneshot
ExecStart={{ atlas_homeserver_root }}/taskchampion-backup.sh
